# ansible/roles/leadgen-api/templates/leadgen-api.container.j2
# Generated by Ansible (Motorcade LeadGen)
# Quadlet file -> systemd unit: {{ leadgen_service_name }}.service

[Container]
ContainerName={{ leadgen_service_name }}
Image={{ leadgen_image_name }}:{{ leadgen_image_tag }}
PublishPort={{ leadgen_host_port }}:{{ leadgen_container_port }}

# Read secrets from host-managed env file
EnvironmentFile=/run/secrets/leadgen.env

# SELinux-safe mounts
Volume={{ leadgen_secret_dir }}:/run/secrets:Z
Volume={{ leadgen_log_root }}:/var/log/leadgen:Z

# Basic hardening
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
